# infra/systemd/orkestai-web.service
[Unit]
Description=ORKESTAI ENGAGE Web Dashboard
After=network.target orkestai-api.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/engage/apps/web
ExecStart=/home/ec2-user/.nvm/versions/node/v22.22.3/bin/node /home/ec2-user/engage/apps/web/node_modules/next/dist/bin/next start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=-/home/ec2-user/engage/.env
Environment="NODE_ENV=production"
Environment="HOSTNAME=0.0.0.0"
Environment="PORT=3000"

[Install]
WantedBy=multi-user.target
