# infra/systemd/orkestai-web.service
[Unit]
Description=ORKESTAI ENGAGE Web Dashboard
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/engage/apps/web/.next/standalone
ExecStart=/home/ec2-user/.nvm/versions/node/v22.22.3/bin/node apps/web/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"
Environment="HOSTNAME=0.0.0.0"
Environment="PORT=3000"
Environment="INTERNAL_API_URL=http://localhost:3001"
Environment="NEXT_PUBLIC_API_URL=http://localhost:3001"

[Install]
WantedBy=multi-user.target
