# infra/systemd/orkestai-api.service
[Unit]
Description=ORKESTAI ENGAGE API Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/engage
ExecStart=/home/ec2-user/.nvm/versions/node/v22.22.3/bin/node apps/api/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"
Environment="DATABASE_URL=postgresql://engage:engage@localhost:5432/engage"
Environment="REDIS_URL=redis://localhost:6379"
Environment="CORS_ORIGINS=https://engage.orkestai.ar,http://localhost:3000"
Environment="API_JWT_SECRET=dev-secret-this-must-be-at-least-32-chars-long-for-production"
Environment="ENCRYPTION_KEY=dev-key-this-must-be-exactly-32-bytes-for-aes-256!"

[Install]
WantedBy=multi-user.target
