# infra/systemd/orkestai-api.service
[Unit]
Description=ORKESTAI ENGAGE API Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/engage
ExecStart=/home/ec2-user/.nvm/versions/node/v22.22.3/bin/node apps/api/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=-/home/ec2-user/engage/.env
Environment="NODE_ENV=production"
Environment="CORS_ORIGINS=https://engage.orkestai.ar,http://localhost:3000"

[Install]
WantedBy=multi-user.target
